# .config/systemd/user/check_battery.timer
[Unit]
Description=Timer to check battery status

[Timer]
OnCalendar=*:0/2
Unit=check_battery.service

[Install]
WantedBy=timers.target
